# meta-scobc-dualboot/recipes-apps/mark-boot-success/files/mark-boot-success.service
[Unit]
Description=Mark successful Linux boot by restoring PMC_MULTI_BOOT register
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/multiboot-tool sd-prev
StandardOutput=journal
StandardError=journal
